# linux-tasks/add-tailscale-route-rule/add-tailscale-route-rule.service
[Unit]
Description=Tailscale routing rule setup
After=network.target network-online.target tailscaled.service
Requires=tailscaled.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/add-tailscale-route-rule.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
