# linux-tasks/check-external-drives-mounted/check-external-drives-mounted.timer
[Unit]
Description=Run external drives check hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
