# linux-tasks/check-mount-points.timer
[Unit]
Description=Run mount points check hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
