# linux-tasks/check-mount-points/check-mount-points.timer
[Unit]
Description=Run checking mount points in every minute

[Timer]
OnCalendar=*:0/1
Persistent=true

[Install]
WantedBy=timers.target
